# deploy/jornal.service
[Unit]
Description=Jornal Node API + Frontend
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/var/www/jornal
Environment=NODE_ENV=production
EnvironmentFile=/var/www/jornal/.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
